# bin/jojoblog.conf
description "Feather Blog Sample App Ubuntu Upstart configuration file"
author      "smk - Steve Krenek"
# Startup (daemonized) script for Ubuntu

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
#respawn
#respawn limit 99 5

#pid file /var/run/blog.pid

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/srv/node/feather/blog"
    chdir $HOME

    exec sudo -u www-data /usr/local/nvm/v0.4.1/bin/node app.js
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script
